# smf/manifest/cainstsvc.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright (c) 2014, Joyent, Inc.
-->

<!-- 
  This manifest is only used by the NPM package for cainstsvc.  It relies on
  several fields being filled in dynamically via the post-activate script
  (these fields are delimited with '@@').  When a better and/or more generic
  mechanism is available (e.g., generic templating of a manifest or dynamically
  modifying a new service with svccfg), this manifest should be changed to use
  it, and unified with its sibling, smartdc-ca-cainstsvc.xml.
  -->

<service_bundle type="manifest" name="smartdc-ca-cainstsvc">
<service name="smartdc/agent/ca/cainstsvc" type="service" version="1">

<dependency name="network" grouping="require_all" restart_on="error" type="service">
	<service_fmri value="svc:/network/physical" />
</dependency>
<dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
	<service_fmri value="svc:/system/filesystem/local" />
</dependency>

<exec_method type="method" name="start" exec="@@CABASE_DIR@@/smf/method/canodesvc" timeout_seconds="30" />
<exec_method type="method" name="stop" exec=":kill" timeout_seconds="30" />

<instance name="@@INSTANCE_NAME@@" enabled="true">
	<property_group name="com.joyent.ca,caconfig" type="application">
		<propval name="path-logs" type="astring" value="@@BASE_DIR@@/log/ca" />
		<propval name="path-node" type="astring" value="@@CABASE_DIR@@/cmd/node" />
		<propval name="path-dist" type="astring" value="@@CABASE_DIR@@" />
		<propval name="run-rt" type="astring" value="true" />
		<propval name="dtrace-libpath" type="astring" value="@@CABASE_DIR@@/lib" />
	</property_group>
</instance>

<stability value='Unstable' />

<template>
	<common_name><loctext xml:lang="C">Cloud Analytics Instrumenter Service</loctext></common_name>
</template>

</service>
</service_bundle>
