# smf/manifest/smartdc-ca-cainstsvc.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type="manifest" name="smartdc-ca-cainstsvc">
<service name="smartdc/ca/cainstsvc" type="service" version="1">

<dependency name="network" grouping="require_all" restart_on="error" type="service">
	<service_fmri value="svc:/network/physical" />
</dependency>
<dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
	<service_fmri value="svc:/system/filesystem/local" />
</dependency>

<exec_method type="method" name="start" exec="/opt/smartdc/ca/smf/method/canodesvc" timeout_seconds="10" />
<exec_method type="method" name="stop" exec=":kill" timeout_seconds="10" />

<instance name="default" enabled="true">
	<property_group name="com.joyent.ca,caconfig" type="application">
		<propval name="path-logs" type="astring" value="/var/log/smartdc/ca" />
		<propval name="path-node" type="astring" value="/opt/smartdc/ca/deps/node/node" />
		<propval name="path-dist" type="astring" value="/opt/smartdc/ca" />
	</property_group>
</instance>

<stability value='Unstable' />

<template>
	<common_name><loctext xml:lang="C">Cloud Analaytics Instrumenter Service</loctext></common_name>
</template>

</service>
</service_bundle>
